# etc/cron.hourly/2snapshot
#!/bin/sh

python3 "/usr/bin/snapshot-home.py"
